# packaging/systemd-user/aperture-gc.service
[Unit]
Description=Trigger Aperture garbage collection

[Service]
Type=oneshot
ExecStart=%h/.nix-profile/bin/aperture trigger gc
