# packaging/docker/rootfs/etc/services.d/traefik/run
#!/command/with-contenv sh
exec 2>&1
exec s6-setuidgid aperture traefik --configFile=/etc/aperture/traefik.yaml
